# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
user=www-data

[program:horizon]
directory=/var/www/html
command=php artisan horizon
user=www-data